# .config/autostart/coolero.desktop
[Desktop Entry]
Name=coolero
Comment=Fan contoller
Exec=/usr/bin/coolero
Terminal=false
Type=Application
StartupNotify=false