# .config/autostart/picom.desktop
[Desktop Entry]
Name=picom
Comment=A compositor for X11
Exec=/usr/bin/picom
Terminal=false
Type=Application
StartupNotify=false