# .config/autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Comment=Removable media auto mounter
Exec=/usr/bin/udiskie -At
Terminal=false
Type=Application
StartupNotify=false