# etc/xdg/autostart/xfce4-notifyd.desktop
[Desktop Entry]
Type=Application
Name=Xfce Notification Daemon
Exec=false
Icon=org.xfce.notification
OnlyShowIn=XFCE;
